# debian/DEBIAN/prerm
#!/bin/bash

echo "Start to uninstall hobot-utils"

update-rc.d hobot-adbd remove

systemctl disable hobot-adbd.service

exit 0

